# os-sim/debian/ossim.postrm
#!/bin/bash
update-rc.d -f ossim remove 2&> /dev/null

