# os-sim/debian/ossim-framework.cron.d
# ossim graphing
0-59/5 * * * * root /usr/share/ossim/mrtg/launch-mrtg

# ossim db backup
#5 0 * * * root /usr/share/ossim/scripts/backupdb.pl

