# os-sim/debian/ossim-agent.conffiles
/etc/ossim/agent/config.xml
/etc/ossim/agent/plugins/arpwatch.xml
/etc/ossim/agent/plugins/apache.xml
/etc/ossim/agent/plugins/camonitor.xml
/etc/ossim/agent/plugins/ca.xml
/etc/ossim/agent/plugins/ciscoids.xml
/etc/ossim/agent/plugins/ciscopix.xml
/etc/ossim/agent/plugins/ciscorouter.xml
/etc/ossim/agent/plugins/fw1.xml
/etc/ossim/agent/plugins/iis.xml
/etc/ossim/agent/plugins/iptables.xml
/etc/ossim/agent/plugins/netgear.xml
/etc/ossim/agent/plugins/ntop.xml
/etc/ossim/agent/plugins/ntsyslog.xml
/etc/ossim/agent/plugins/opennms.xml
/etc/ossim/agent/plugins/osiris.xml
/etc/ossim/agent/plugins/p0f.xml
/etc/ossim/agent/plugins/pads.xml
/etc/ossim/agent/plugins/prelude.xml
/etc/ossim/agent/plugins/realsecure.xml
/etc/ossim/agent/plugins/rrd.xml
/etc/ossim/agent/plugins/snarewindows.xml
/etc/ossim/agent/plugins/snort.xml
/etc/ossim/agent/plugins/syslog.xml
/etc/ossim/agent/plugins/tcptrack.xml
/etc/init.d/ossim-agent
