# os-sim/debian/ossim.postinst
#!/bin/bash
update-rc.d ossim defaults 40 2&> /dev/null
/etc/init.d/ossim restart

